# systemd/Modello_Homescreen-launchpad-ready.path
# Description=Path activation for Modello_Homescreen

[Path]
PathExists=/run/user/%U/wayland-0
Unit=Modello_Homescreen.service

[Install]
WantedBy=weston.target
